# files/etc/s6/consul-template/run
#!/bin/sh

exec /usr/local/bin/consul-template -consul $CONSUL_HOST -config /etc/consul/template.cfg
